# supervision/supervisor/komorebi.ini
[program:komorebi]
command=/usr/local/bin/python3 /home/freebsd/transitional/komorebi/komorebi/run.py --unix /var/run/www/komorebi.socket
process_name=%(program_name)s
directory=/home/freebsd/transitional/komorebi
user=www
redirect_stderr=true
stdout_logfile=/var/log/www/talideon.com.komorebi.log
stderr_logfile=/dev/null
environment=KOMOREBI_DB_PATH="/var/db/komorebi/db.sqlite",PYTHONPATH="/home/freebsd/transitional/komorebi"
