# supervision/komorebi.ini
[program:komorebi]
command=/usr/bin/local/python3 -m komorebi
process_name=%(program_name)s
directory=/path/to/app
user=www
redirect_stderr=true
stdout_logfile=/var/log/www/komorebi.log
stderr_logfile=/dev/null
environment=KOMOREBI_SETTINGS="/path/to/app/komorebi.cfg",PYTHONPATH="/path/to/app"
